# samples/sample-image/with-s6-overlay/s6/s6-overlay/s6-rc.d/app/run
#!/bin/bash
exec 2>&1
cd /
exec uvicorn main:app --host '*' --port 80
